# service/lxdrunner.openrc
#!/sbin/openrc-run

name="LXDRunner"
description="Emphemeral selfhosted runners"
command="/home/app/.local/bin/lxdrunner"
pidfile="/var/run/$SVCNAME.pid"
command_background=true
command_user="app:app"

depend() {
	need logger
}

